# SCS/student_setup/notebooks/scripts/crontab
*/2 * * * * /home/pi/notebooks/scripts/singleimage.py
